# unbindgpu/DEBIAN/postinst
#!/bin/sh

# Enable the service
systemctl enable unbindgpu.service || exit
