# unbindgpu/DEBIAN/prerm
#!/bin/sh

# Disable the service
systemctl disable unbindgpu.service || exit
